# templates/ath11k/ath11k_pci.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="ath11k_pci", RUN+="/sbin/ath11k_disable_power_management.sh"
